# pitlab.conf
description "Upstart script to run pitlab as a service"

env BIN=pitlab
env LOG_FILE=/var/log/pitlab.log
env CONF_FILE=/home/git/pitlab.json
env PORT=3000
env RUN_AS="pitlab"        # Upstart can only be run nicely as root, need to drop privileges

# Start service on startup, stop on shutdown
start on runlevel [2345]
stop on runlevel [016]

# Respawn in case of a crash, with default parameters
respawn

script
        # Make sure logfile exists and can be written by the user we drop privileges to
        touch $LOG_FILE
        chown $RUN_AS:$RUN_AS $LOG_FILE

        PITLAB_PORT=$PORT su -s /bin/sh -c 'exec "$0" "$@"' $RUN_AS -- $BIN $CONF_FILE >> $LOG_FILE 2>&1
end script


post-start script
        echo "===== App restarted =====" >> $LOG_FILE
end script
